# .config/autostart/nitrogen.desktop
[Desktop Entry]
Type=Application
Name=Nitrogen Wallpaper
Exec=nitrogen --restore
