# .config/autostart/nm-applet.desktop
[Desktop Entry]
Type=Application
Name=NetworkManager applet
Exec=nm-applet
