# .config/autostart/picom.desktop
[Desktop Entry]
Type=Application
Name=Picom Compositor
Exec=picom -f
