# .config/autostart/polybar.desktop
[Desktop Entry]
Type=Application
Name=Polybar with sys tray
Exec=/home/pengyul/.config/polybar/launch.sh --material
