# .config/autostart/xbindkeys.desktop
[Desktop Entry]
Type=Application
Name=Key bindings
Exec=xbindkeys
